# services/znc/run
#!/bin/bash

umask 002
exec /sbin/setuser abc /usr/local/bin/znc -d /config
