# defaults/supervisord.conf
[supervisord]
nodaemon=true
[program:znc]
user=abc
command=/usr/local/bin/znc -d /config

