# root/etc/services.d/znc/run
#!/usr/bin/with-contenv sh

umask 002

exec \
	s6-setuidgid abc /usr/bin/znc -d /config \
	--foreground --debug "$@"
